# container-files/etc/supervisor.d/dropboxd.conf
[program:dropboxd]
command=/opt/dropbox/dropbox-dist/dropboxd
autostart=true
stderr_logfile = NONE
stdout_logfile = /var/log/secure
redirect_stderr = true
